# templates/default/lib/systemd/system/unicorn-systemd.service.erb
[Unit]
Description=<%= @description %>

[Service]
Environment="GMAIL_DOMAIN=gmail.com"
Environment="GMAIL_USERNAME=<%= @gmail_username %>"
Environment="GMAIL_PASSWORD=<%= @gmail_password %>"

WorkingDirectory=<%= @cwd %>
ExecStart=<%= @exec_start %>
ExecStop=<%= @exec_stop %>
Restart=on-failure

[Install]
WantedBy=multi-user.target

